# supervisor.ini
[supervisord]
nodaemon=true
[program:pythonstudy]
directory=/src
command=python /src/manage.py runserver 0.0.0.0:8000
stopsignal=QUIT
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/%(program_name)s.log
[program:celery]
directory=/src
command=celery -A pythonstudy worker -l info
stopsignal=QUIT
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/%(program_name)s.log